# roles/xui/templates/3x-ui-backup.service.j2
[Unit]
Description=Backup 3x-ui configuration and certificates
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/3x-ui-backup.sh