# roles/xui/templates/3xui-panel-access-guard.timer.j2
[Unit]
Description=Периодическая проверка времени прямого доступа панели 3x-ui

[Timer]
OnBootSec=2m
OnUnitActiveSec={{ xui_panel_access_guard_interval }}
AccuracySec=30s
Unit=3xui-panel-access-guard.service
Persistent=true

[Install]
WantedBy=timers.target
