# roles/sni_manager/templates/3x-ui-manager-rotation.timer.j2
[Unit]
Description=Планировщик ротации Reality SNI

[Timer]
OnCalendar={{ xui_manager_rotation_on_calendar }}
Persistent=true
Unit={{ xui_manager_rotation_service_name }}.service

[Install]
WantedBy=timers.target